# com.vehicle-vitals.runner-token-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vehicle-vitals.runner-token-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/marknelson/Circus/Repositories/vehicle-vitals/token-refresh.sh</string>
      <string>health_check</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer> <!-- Run every hour -->

    <key>StandardOutPath</key>
    <string>/tmp/runner-token-refresh.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/runner-token-refresh.err</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <false />
  </dict>
</plist>